# install-server.conf
start on runlevel 2

exec LiRE/install_server.py 56000 
