# lire-install-server.conf
start on runlevel 2

exec /usr/local/bin/install_server.py 56000 
